# src/test/resources/docker/cluster/redis-cluster/conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:redis-50000]
command=/usr/bin/redis-server /opt/redis/50000/redis-50000.conf
stderr_logfile = /var/log/redis-50000-stderr.log
stdout_logfile = /var/log/redis-50000-stdout.log
priority=100
directory=/opt/redis/50001

[program:redis-50001]
command=/usr/bin/redis-server /opt/redis/50001/redis-50001.conf
stderr_logfile = /var/log/redis-50001-stderr.log
stdout_logfile = /var/log/redis-50001-stdout.log
priority=100
directory=/opt/redis/50001

[program:redis-50002]
command=/usr/bin/redis-server /opt/redis/50002/redis-50002.conf
stderr_logfile = /var/log/redis-50002-stderr.log
stdout_logfile = /var/log/redis-50002-stdout.log
priority=100
directory=/opt/redis/50002

[program:redis-50003]
command=/usr/bin/redis-server /opt/redis/50003/redis-50003.conf
stderr_logfile = /var/log/redis-50003-stderr.log
stdout_logfile = /var/log/redis-50003-stdout.log
priority=100
directory=/opt/redis/50003

[program:redis-50004]
command=/usr/bin/redis-server /opt/redis/50004/redis-50004.conf
stderr_logfile = /var/log/redis-50004-stderr.log
stdout_logfile = /var/log/redis-50004-stdout.log
priority=100
directory=/opt/redis/50004

[program:redis-50005]
command=/usr/bin/redis-server /opt/redis/50005/redis-50005.conf
stderr_logfile = /var/log/redis-50005-stderr.log
stdout_logfile = /var/log/redis-50005-stdout.log
priority=100
directory=/opt/redis/50005

[program:create-cluster]
command=/usr/bin/redis-cli --cluster create 127.0.0.1:50000 127.0.0.1:50001 127.0.0.1:50002 127.0.0.1:50003 127.0.0.1:50004 127.0.0.1:50005 --cluster-replicas 1 --cluster-yes
stderr_logfile = /var/log/redis-cluster-stderr.log
stdout_logfile = /var/log/redis-cluster-stdout.log








