# src/test/resources/docker/sentinel/redis-sentinel/conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket


[program:redis-50000]
command=/usr/bin/redis-server /opt/redis/redis-50000.conf
stderr_logfile = /var/log/redis-50000-stderr.log
stdout_logfile = /var/log/redis-50000-stdout.log

[program:redis-50001]
command=/usr/bin/redis-server /opt/redis/redis-50001.conf
stderr_logfile = /var/log/redis-50001-stderr.log
stdout_logfile = /var/log/redis-50001-stdout.log

[program:redis-50002]
command=/usr/bin/redis-server /opt/redis/redis-50002.conf
stderr_logfile = /var/log/redis-50002-stderr.log
stdout_logfile = /var/log/redis-50002-stdout.log

[program:redis-sentinel-51000]
command=/usr/bin/redis-sentinel /opt/redis/redis-sentinel-51000.conf
stderr_logfile = /var/log/redis-sentinel-51000-stderr.log
stdout_logfile = /var/log/redis-sentinel-51000-stdout.log

[program:redis-sentinel-51001]
command=/usr/bin/redis-sentinel /opt/redis/redis-sentinel-51001.conf
stderr_logfile = /var/log/redis-sentinel-51001-stderr.log
stdout_logfile = /var/log/redis-sentinel-51001-stdout.log

[program:redis-sentinel-51002]
command=/usr/bin/redis-sentinel /opt/redis/redis-sentinel-51002.conf
stderr_logfile = /var/log/redis-sentinel-51002-stderr.log
stdout_logfile = /var/log/redis-sentinel-51002-stdout.log
